# deploy/aerocrawl.service
[Unit]
Description=Aerocrawl - Internal Scraping Service
After=network.target redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/aerocrawl
EnvironmentFile=/opt/aerocrawl/.env
ExecStart=/opt/aerocrawl/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8001 --workers 1
Restart=always
RestartSec=5
Environment=TZ=Asia/Kolkata
ReadWritePaths=/opt/aerocrawl/data /tmp /root/.cache/ms-playwright

[Install]
WantedBy=multi-user.target
